# autostart/rdp-credentials.desktop
[Desktop Entry]
Type=Application
Name=RDP Credentials Setup
Comment=Set GNOME Remote Desktop credentials on login
Exec=bash -c "sleep 5 && grdctl rdp enable && grdctl rdp set-credentials yaniv testpass123 && grdctl rdp disable-view-only"
Hidden=false
NoDisplay=true
X-GNOME-Autostart-enabled=true
X-GNOME-Autostart-Delay=5
