# utilities/chrome-permission-guard.service
[Unit]
Description=AutomateLinux Chrome Permission Guard
After=network.target

[Service]
ExecStart=/opt/automateLinux/utilities/chrome-permission-guard.sh
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
