# services/system/dashboard-bridge.service
[Unit]
Description=Dashboard Bridge (daemon-to-frontend bridge)
After=daemon.service
Wants=daemon.service
StartLimitIntervalSec=300
StartLimitBurst=20

[Service]
Type=simple
WorkingDirectory=/opt/automateLinux/dashboard
ExecStartPre=-/bin/bash -c 'fuser -k 3501/tcp 2>/dev/null || true; sleep 0.5'
ExecStart=/usr/bin/node bridge.cjs
Restart=on-failure
RestartSec=5s
User=root
Group=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
