# services/system/dashboard-dev.service
[Unit]
Description=Dashboard Dev (Vite)
After=dashboard-bridge.service
Wants=dashboard-bridge.service

[Service]
Type=simple
WorkingDirectory=/opt/automateLinux/dashboard
ExecStart=/usr/bin/npx vite --port 3007 --host
Restart=on-failure
RestartSec=5s
User=root
Group=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
