# services/system/loom-server-dev.service
[Unit]
Description=Loom Server Recording Backend (dev)
After=network.target daemon.service
Wants=daemon.service

[Service]
Type=simple
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=WAYLAND_DISPLAY=wayland-0
Environment=XDG_SESSION_TYPE=wayland
Environment=XDG_CURRENT_DESKTOP=Unity
Environment=XDG_SESSION_DESKTOP=ubuntu
Environment=XDG_MENU_PREFIX=gnome-
WorkingDirectory=/opt/automateLinux/extraApps/loom/server/build
ExecStart=/opt/automateLinux/extraApps/loom/server/build/loom-server 3505 dev
Restart=always
RestartSec=5s
User=yaniv
Group=yaniv
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
