# services/system/loom-client-dev.service
[Unit]
Description=Loom Client Dev (Vite)
After=loom-server-dev.service
Wants=loom-server-dev.service

[Service]
Type=simple
WorkingDirectory=/opt/automateLinux/extraApps/loom/client
ExecStart=/usr/bin/npx vite --port 3005 --host
Restart=on-failure
RestartSec=5s
User=root
Group=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
