# services/system/loom-client-prod.service
[Unit]
Description=Loom Client Prod (Vite)
After=loom-server-prod.service
Wants=loom-server-prod.service

[Service]
Type=simple
WorkingDirectory=/opt/prod/loom/client
ExecStart=/usr/bin/npx vite --port 3004 --host
Restart=on-failure
RestartSec=5s
User=root
Group=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
